# contrib/init/bitcoind.conf
description "ReBitcoin Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env REBITCOIND_BIN="/usr/bin/rebitcoind"
env REBITCOIND_USER="rebitcoin"
env REBITCOIND_GROUP="rebitcoin"
env REBITCOIND_PIDDIR="/var/run/rebitcoind"
# upstart can't handle variables constructed with other variables
env REBITCOIND_PIDFILE="/var/run/rebitcoind/rebitcoind.pid"
env REBITCOIND_CONFIGFILE="/etc/rebitcoin/rebitcoin.conf"
env REBITCOIND_DATADIR="/var/lib/rebitcoind"

expect fork

respawn
respawn limit 5 120
kill timeout 600

pre-start script
    # this will catch non-existent config files
    # rebitcoind will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$REBITCOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run rebitcoind."
        echo "The setting must appear in $REBITCOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"ReBitcoin Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$REBITCOIND_PIDDIR"
    chmod 0755 "$REBITCOIND_PIDDIR"
    chown $REBITCOIND_USER:$REBITCOIND_GROUP "$REBITCOIND_PIDDIR"
    chown $REBITCOIND_USER:$REBITCOIND_GROUP "$REBITCOIND_CONFIGFILE"
    chmod 0660 "$REBITCOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$REBITCOIND_PIDFILE" \
    --chuid $REBITCOIND_USER:$REBITCOIND_GROUP \
    --exec "$REBITCOIND_BIN" \
    -- \
    -pid="$REBITCOIND_PIDFILE" \
    -conf="$REBITCOIND_CONFIGFILE" \
    -datadir="$REBITCOIND_DATADIR" \
    -disablewallet \
    -daemon

